# build/systemd/kvm-agent.service
[Unit]
Description=KVM Agent
After=network.target

[Service]
ExecStart=/usr/bin/kvm-agent -c /etc/kvm-agent/config.yaml
ExecStop=/bin/kill -TERM $MAINPID
Type=simple
Restart=no
PIDFile=/var/run/kvm-agent

[Install]
WantedBy=multi-user.target